# ngcp-config/templates/etc/cron.d/ngcp-credit-warning.tt2
13 6 * * *  root  if /usr/sbin/ngcp-check_active -q; then /usr/sbin/ngcp-credit-warning; fi
